# deploy/review-me.service
[Unit]
Description=review-me PR review bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=review-me
Group=review-me
WorkingDirectory=/opt/review-me
EnvironmentFile=/etc/review-me/review-me.env
ExecStart=/usr/local/bin/bun run /opt/review-me/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Harden the service — adjust if your deployment needs broader access.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
